# servers/electris_cron.conf
description "Pseudo cron process for fast electris updates"

start on runlevel [2345]
stop on runlevel [!2345]

respawn

setuid ubuntu
setgid ubuntu
env HOME=/home/ubuntu
env DEPLOYMENT_TARGET=staging

script
    chdir /home/ubuntu/apps/electris/repository
    exec /home/ubuntu/apps/electris/virtualenv/bin/fab $DEPLOYMENT_TARGET deployment_cron
end script